# hello-app.service
[Unit]
Description=Hello App Service

[Service]
ExecStart=/usr/bin/hello-app

[Install]
WantedBy=multi-user.target
